# debian/android-propd.property-service.socket
[Unit]
Description=Android property service socket

[Socket]
ListenStream=/run/property_service

[Install]
WantedBy=sockets.target
